# fab/services/templates/supervisor_submission_reprocessing_queue.conf
[program:{{ project }}-{{ environment }}-submission_reprocessing_queue]
directory={{ code_current }}/
command={{ virtualenv_current }}/bin/python manage.py run_submission_reprocessing_queue
user={{ sudo_user }}
numprocs=1
autostart=true
autorestart=true
stdout_logfile={{ log_dir }}/submission_reprocessing_queue.log
redirect_stderr=true
stderr_logfile={{ log_dir }}/submission_reprocessing_queue.error.log
startsecs=0
stopwaitsecs=10
