# ansible/roles/supervisor/templates/supervisord.upstart.conf.j2
description     "supervisord"

start on runlevel [2345]
stop on runlevel [!2345]

respawn

exec {{ which_supervisord.stdout }} --nodaemon --configuration /etc/supervisord.conf
