# src/commcare_cloud/ansible/roles/supervisor/templates/supervisor_task.conf.j2
;
; {{ ansible_managed }}
;
[program:{{ item.name }}]
{% for key, val in item.config.items() %}
{{ key }}={{ val }}
{% endfor %}
