# fab/fab/services/templates/supervisor_queue_schedule_instances.conf
[program:{{ project }}-{{ environment }}-queue_schedule_instances]
directory={{ code_current }}/
command={{ virtualenv_current }}/bin/python manage.py queue_schedule_instances
user={{ sudo_user }}
numprocs=1
autostart=true
autorestart=true
stdout_logfile={{ log_dir }}/queue_schedule_instances.log
redirect_stderr=true
stderr_logfile={{ log_dir }}/queue_schedule_instances.error.log
startsecs=0
stopwaitsecs=10
