# src/commcare_cloud/ansible/roles/commcarehq/templates/supervisor_websockets.conf
[program:{{ project }}-{{ deploy_env }}-websockets]
directory={{ code_home }}
command={% if PY3_RUN_WEBSOCKETS %}{{ py3_virtualenv_home }}{% else %}{{ virtualenv_home }}{% endif %}/bin/uwsgi
    --virtualenv {% if PY3_RUN_WEBSOCKETS %}{{ py3_virtualenv_home }}{% else %}{{ virtualenv_home }}{% endif %}
    --http-socket {{ www_home }}/web.socket
    --gevent 1000
    --http-websockets
    --workers=2
    --master
    --module deployment.websocket_wsgi
    --chmod-socket=664
user={{ cchq_user }}
numprocs=1
autostart=true
autorestart=true
stdout_logfile={{ log_home }}/websockets.log
redirect_stderr=true
stderr_logfile={{ log_home }}/websockets.error.log
startsecs=0
stopwaitsecs=10
